# src/devsn2mac.service
[Unit]
Description=Configure eth0 MAC address from DEVSN

[Service]
Type=oneshot
ExecStart=/usr/sbin/devsn2mac

[Install]
WantedBy=network.target
