# src/setmac.service
[Unit]
Description=Configure eth0 MAC address from NVRAM

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/setmac

[Install]
WantedBy=network.target
